# supervisor-conf.d/camera_mod.conf
[program:camera_mod]
command=/home/gserv/venv/bin/python -m gserv.CameraModule
directory=/home/gserv/GarageServerPi
user=gserv
autostart=true
autorestart=unexpected
startretries=3
