# supervisor-conf.d/controller_mod.conf
[program:controller_mod]
command=/home/gserv/venv/bin/python -m gserv.ControllerModule
directory=/home/gserv/GarageServerPi
user=gserv
autostart=true
autorestart=unexpected
startretries=3
