# supervisor-conf.d/gpioinput_mod.conf
[program:gpioinput_mod]
command=/home/gserv/venv/bin/python -m gserv.GPIOInputModule
directory=/home/gserv/GarageServerPi
user=root
autostart=true
autorestart=unexpected
startretries=3
