# supervisor-conf.d/gpiooutput_mod.conf
[program:gpiooutput_mod]
command=/home/gserv/venv/bin/python -m gserv.GPIOOutputModule
directory=/home/gserv/GarageServerPi
user=root
autostart=true
autorestart=unexpected
startretries=3
