# supervisor-conf.d/heartbeat_mod.conf
[program:heartbeat_mod]
command=/home/gserv/venv/bin/python -m gserv.HeartbeatModule
directory=/home/gserv/GarageServerPi
user=gserv
autostart=true
autorestart=unexpected
startretries=3
