# supervisor-conf.d/led_module.conf
[program:led_mod]
command=/home/gserv/venv/bin/python -m gserv.LEDSModule
directory=/home/gserv/GarageServerPi
user=root
autostart=true
autorestart=unexpected
startretries=3
