# supervisor-conf.d/thingspeak_mod.conf
[program:thingspeak_mod]
command=/home/gserv/venv/bin/python -m gserv.ThingspeakModule
directory=/home/gserv/GarageServerPi
user=gserv
autostart=true
autorestart=unexpected
startretries=3
